# system/systemd/user/bots.timer
[Unit]
Description=Post co2 ppm to slack

[Timer]
OnUnitActiveSec=60s

[Install]
WantedBy=timers.target
